# infra/systemd/imob-crawl-seed.service
[Unit]
Description=ImobIntel seed-discover (enqueue discover jobs for all listing sources)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'curl -sf -H "Authorization: Bearer $CRON_SECRET" "${CRON_BASE_URL:-https://imobintel.ro}/api/cron/seed-discover"'
EnvironmentFile=/home/admin/apps/imob/infra/.env
TimeoutStartSec=120
