# infra/systemd/imob-dedup-tick.service
[Unit]
Description=ImobIntel dedup-tick (group duplicate listings)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -sf -H "Authorization: Bearer ${CRON_SECRET}" http://localhost:3000/api/cron/dedup-tick
EnvironmentFile=/home/admin/apps/imob/infra/.env
TimeoutStartSec=90
