# infra/systemd/imob-avm-train.service
[Unit]
Description=ImobIntel AVM model training (ridge regression on historical data)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -sf -H "Authorization: Bearer ${CRON_SECRET}" http://localhost:3000/api/cron/avm-train
EnvironmentFile=/home/admin/apps/imob/infra/.env
TimeoutStartSec=300
