# infra/systemd/imob-avm-train.service
[Unit]
Description=ImobIntel AVM model training (ridge regression on historical data)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'curl -sf -H "Authorization: Bearer $CRON_SECRET" "${CRON_BASE_URL:-https://imobintel.ro}/api/cron/avm-train"'
EnvironmentFile=/home/admin/apps/imob/infra/.env
TimeoutStartSec=300
